# ecs-ec2-spot-fleet-mustacheme/mustacheme/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx

[program:uswgi]
command=/var/www/mustacheme/venv/bin/uwsgi --ini /var/www/mustacheme/nginx/uwsgi.ini

